# openwrt/shawnwrt-ota/files/etc/uci-defaults/99-shawnwrt-ota-finalize
#!/bin/sh

/usr/bin/shawnwrt-ota finalize-pending >/dev/null 2>&1 || true

exit 0
